# packages/yacht/captain.service
[Unit]
Description=captain service
After=network.target

[Service]
ExecStartPre=/bin/bash -c 'git clone git@github.com:timmyg/captain.git'
ExecStartPre=/bin/bash -c 'cd captain/packages/yacht'
ExecStartPre=/bin/bash -c 'git pull origin master'
ExecStartPre=/bin/bash -c 'yarn install'
ExecStart=/bin/bash -c 'npm start'
WorkingDirectory=/home/pi
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target

sudo systemctl enable myfirst